# minecraft-server-workaround-7069-after.service
[Unit]
Description=Minecraft server auxillary unit to workaround bug systemd#7069 (After)
BindsTo=minecraft-server.service
After=minecraft-server.service
JoinsNamespaceOf=minecraft-server.service
RefuseManualStart=yes

[Service]
User=minecraft-server
Group=minecraft-server
RemainAfterExit=yes
ExecStop=/bin/sh -c 'echo stop >/run/minecraft-server-stdin.fifo && while x="$(LC_ALL=C systemctl is-active minecraft-server.service)"; { [ "$x" = active ] || [ "$x" = deactivating ]; }; do sleep 1; done'
TimeoutStopSec=1 minute

# Hardening options
CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=full
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
PrivateTmp=yes
LockPersonality=yes
