# minecraft-server@.socket
[Unit]
Description=Minecraft server game data socket (%i)

[Socket]
ListenStream=25565

[Install]
WantedBy=sockets.target
